# django/supervisord/conf.supervisord/uwsgi.conf
[program:uwsgi]
# Project directory
directory=/code
# Application
command=uwsgi --ini /uwsgi.ini
user=root
redirect_stderr=true
autostart=true
autorestart=true
