# nginx/supervisord/conf.supervisord/cron.conf
[program:cron]
command=/usr/sbin/crond -f
stdout_syslog=true
stderr_syslog=true
autorestart=false
startentries=0
