# nginx/supervisord/conf.supervisord/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stdout_syslog=true
stderr_syslog=true
autorestart=false
startentries=0
